# cookbooks/basethin/templates/default/thin.service.erb
# This file was generated by the Chef 'basethin' cookbook (updated <%= Time.now.utc.strftime('%d %b %Y %H:%M') %>)
[Unit]
Description=A fast and very simple Ruby web server
After=syslog.target network.target

[Service]
Type=simple
User=<%= @thin_user %>
Group=<%= @thin_user %>
ExecStart=/bin/thin start --all <%= @thin_config_dir %>
ExecReload=/bin/thin restart --all <%= @thin_config_dir %>
ExecStop=/bin/thin stop --all <%= @thin_config_dir %>

[Install]
WantedBy=multi-user.target
